# gesture-daemon.service
[Unit]
Description=Touchpad Gesture Daemon
Documentation=man:gesture-daemon(1)
After=multi-user.target graphical-session.target
Wants=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/gesture-daemon
Restart=always
RestartSec=5
User=root
Group=input

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=gesture-daemon
SyslogLevel=info
LogLevelMax=info

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=no
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Device access permissions
DeviceAllow=/dev/hidraw1 r
DeviceAllow=/dev/uinput rw
DevicePolicy=closed

# File system access
ReadWritePaths=/dev/uinput
ReadOnlyPaths=/dev/hidraw1

# Faster shutdown settings
TimeoutStopSec=1          # Kill after 1 second instead of default 90, this is due to blocking nature of read(), which will prevent both booting and stopping the service until a touchpad input is registered and unblocking the read() due to data going to the buffer
KillMode=mixed            # Send SIGTERM, then SIGKILL if needed
SendSIGKILL=yes          # Force kill if graceful shutdown fails

[Install]
WantedBy=multi-user.target
